# deploy/engram.service
[Unit]
Description=Engram core (socket-activated)
Requires=engram.socket
After=engram.socket

[Service]
# Type=notify: engramd inherits the listening fd from engram.socket (LISTEN_FDS) and signals
# readiness via sd_notify - it never binds the port itself (that would EADDRINUSE on the first
# request under socket activation).
Type=notify
ExecStart=/usr/local/bin/engramd
Environment=ENGRAM_HOME=/var/lib/engram
Environment=ENGRAM_IDLE_SECS=120
# To expose Engram on the network, change engram.socket's ListenStream and set a token here:
# Environment=ENGRAM_API_TOKEN=change-me
# Hardening - minimal privileges for a self-modifying agent.
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
StateDirectory=engram

[Install]
WantedBy=multi-user.target
